# personal_config/nas_stuff/mnt-MyStuff.automount
[Unit]
Description=Automount CIFS Remote Tank Automount
Requires=network-online.target
After=network-online.service

[Automount]
Where=/mnt/MyStuff

[Install]
WantedBy=multi-user.target